# codearena_uwsgi.conf
description "uWSGI application server"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid www-data

exec /var/codearena/venv/bin/uwsgi --ini /var/codearena/codearena_uwsgi.ini

